# deploy/projectplanner-mcp.service
# /etc/systemd/system/projectplanner-mcp.service
[Unit]
Description=ProjectPlanner MCP server (:8111, localhost)
After=network-online.target projectplanner-gateway.service
Wants=network-online.target

[Service]
Slice=projectplanner-interactive.slice
Type=simple
User=ubuntu
WorkingDirectory=/opt/projectplanner
EnvironmentFile=/opt/projectplanner/.env
# Production MCP writes must fail closed even if .env is missing PM_AUTH_MODE.
Environment=PM_AUTH_MODE=required
ExecStart=/opt/projectplanner/.venv/bin/python /opt/projectplanner/mcp_server.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
