# deploy/projectplanner-monitors.timer
# /etc/systemd/system/projectplanner-monitors.timer
# Evaluates durable Switchboard coordination monitors every minute.
# To enable:
#   sudo cp projectplanner-monitors.{service,timer} /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now projectplanner-monitors.timer
[Unit]
Description=ProjectPlanner Switchboard monitor sweep (every minute)

[Timer]
OnCalendar=*:0/1
Persistent=true
AccuracySec=10s

[Install]
WantedBy=timers.target
