# deploy/projectplanner-reconcile.timer
# /etc/systemd/system/projectplanner-reconcile.timer
# Runs all-project reconcile alerts every 15 minutes.
# To enable:
#   sudo cp projectplanner-reconcile.{service,timer} /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now projectplanner-reconcile.timer
[Unit]
Description=ProjectPlanner all-project reconcile alerts (every 15 minutes)

[Timer]
OnCalendar=*:0/15
RandomizedDelaySec=3m
Persistent=true
AccuracySec=30s

[Install]
WantedBy=timers.target
