# deploy/projectplanner.service
[Unit]
Description=ProjectPlanner web app (:8110, localhost)
After=network-online.target projectplanner-gateway.service
Wants=network-online.target

[Service]
Slice=projectplanner-interactive.slice
Type=simple
User=ubuntu
WorkingDirectory=/opt/projectplanner
EnvironmentFile=/opt/projectplanner/.env
# Production web writes must fail closed even if .env is missing PM_AUTH_MODE.
Environment=PM_AUTH_MODE=required
ExecStart=/opt/projectplanner/.venv/bin/uvicorn app:app --host 127.0.0.1 --port 8110
Restart=always
RestartSec=3
TimeoutStartSec=60
TimeoutStopSec=15

[Install]
WantedBy=multi-user.target
