# deploy/projectplanner-gateway.service
[Unit]
Description=ProjectPlanner LiteLLM gateway (:8095, localhost)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/projectplanner
EnvironmentFile=/opt/projectplanner/.env
ExecStart=/opt/projectplanner/.venv/bin/litellm --config /opt/projectplanner/deploy/gateway/config.yaml --host 127.0.0.1 --port 8095
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
