# deploy/projectplanner-agent-host.service
# /etc/systemd/system/projectplanner-agent-host.service
# Safe default Agent Host: consumes message-only wake intents for Switchboard.
# It intentionally advertises a sentinel lane so lane-scoped work wakes do not
# dispatch into this service until a separate work-capable host is configured.
[Unit]
Description=ProjectPlanner Switchboard Agent Host (message-only wakes)
After=network-online.target projectplanner.service projectplanner-mcp.service
Wants=network-online.target

[Service]
Slice=projectplanner-interactive.slice
Type=simple
# HARDEN-55: dedicated non-login service account (not the general `ubuntu` login).
User=projectplanner
Group=projectplanner
WorkingDirectory=/opt/projectplanner
EnvironmentFile=/opt/projectplanner/.env
Environment=PM_BASE=https://plan.taikunai.com
Environment=PM_PROJECT=switchboard
Environment=PM_HOST_ID=host/plan-vm-message-wake
Environment=PM_RUNTIME=claude-code
Environment=PM_HOST_LANES=__MESSAGE_ONLY__
Environment=PM_AGENT_HOST_ALLOW_WORK=0
Environment=PM_AGENT_HOST_ALLOW_GLOBAL_CLAIM=0
Environment=PM_HOST_MAX_SESSIONS=1
Environment=PM_REPO_ROOT=/opt/projectplanner
Environment=PM_RUNNER_DIR=/var/lib/projectplanner/runner
Environment=PM_AGENT_HOST_INBOX_IDLE_SECONDS=6
Environment=PYTHONDONTWRITEBYTECODE=1
ExecStart=/opt/projectplanner/.venv/bin/python /opt/projectplanner/adapters/agent_host.py --interval 10
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
# HARDEN-55: runtime least-privilege sandbox. Code tree + venv are root-owned and
# read-only to this service; the only writable path is the data dir (+ journal).
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK
ReadWritePaths=/var/lib/projectplanner

[Install]
WantedBy=multi-user.target
