# deploy/projectplanner-autodeploy.service
# /etc/systemd/system/projectplanner-autodeploy.service
# BUG-114: one-shot auto-deploy trigger. Fetches canonical master and, when the
# running checkout is behind, runs deploy/redeploy.sh (pull -> deps -> units ->
# restart -> Caddy fail-closed -> runtime proof + rollback). Refreshes the
# /health/version staleness signal on every tick. Triggered by the timer below.
#
# WHY THIS UNIT RUNS AS root AND IS NOT SANDBOXED like the batch units
# (reconcile/narrate/etc.): its entire job is to mutate the running system —
# `git pull` on the root-owned code tree (HARDEN-55), `systemctl restart` of the
# web tier, a Caddy reload, AND writing the deploy-state signal into the
# projectplanner-owned data dir (/var/lib/projectplanner, mode 750). No single
# non-root account can do all three: `ubuntu` has sudo but cannot write the
# projectplanner-owned data dir, and `projectplanner` owns the data dir but has
# no sudo. root can, and the 0644 state file it writes is readable by the
# projectplanner web app. So NoNewPrivileges / ProtectSystem=strict cannot apply
# here; redeploy.sh still confines its steps and carries its own rollback guard.
[Unit]
Description=ProjectPlanner Switchboard auto-deploy canonical master to prod
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
Nice=5
User=root
Group=root
WorkingDirectory=/opt/projectplanner
EnvironmentFile=/opt/projectplanner/.env
Environment=PLAN_ROOT=/opt/projectplanner
Environment=PYTHONDONTWRITEBYTECODE=1
# Pin the staleness state file to the service-writable data dir. This MUST match
# what the web app resolves so GET /health/version reads the same file the timer
# writes — set the identical PM_DEPLOY_STATE_FILE in /opt/projectplanner/.env
# (both services load it via EnvironmentFile). An explicit value here also keeps
# the timer off the root-owned repo-root fallback it could never write.
Environment=PM_DEPLOY_STATE_FILE=/var/lib/projectplanner/deploy-state.json
# The script locks the whole cycle (atomic mkdir), so an overlapping tick or a
# tick landing on top of a long redeploy exits 0 without a second deploy.
ExecStart=/usr/bin/env bash /opt/projectplanner/deploy/auto_deploy.sh
# A failed redeploy has already rolled the edge back inside redeploy.sh; surface
# the non-zero result to the journal but never let systemd retry-loop a deploy.
TimeoutStartSec=1200
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
