# deploy/projectplanner-monitors.service
# /etc/systemd/system/projectplanner-monitors.service
# One-shot service: evaluate durable Switchboard coordination monitors.
# Triggered by projectplanner-monitors.timer.
[Unit]
Description=ProjectPlanner Switchboard monitor sweep
After=network.target

[Service]
Slice=projectplanner-batch.slice
Type=oneshot
Nice=10
IOSchedulingClass=idle
MemoryMax=256M
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/projectplanner
EnvironmentFile=/opt/projectplanner/.env
ExecStart=/opt/projectplanner/.venv/bin/python /opt/projectplanner/jobs.py sweep_monitors
StandardOutput=journal
StandardError=journal
