# deploy/projectplanner-inbox.service
# /etc/systemd/system/projectplanner-inbox.service
# Oneshot: poll the Live Inbox mailbox (IMAP) and queue triaged messages. No-op until
# PM_IMAP_* is configured. Triggered by projectplanner-inbox.timer (enable it once the
# mailbox + app password exist).
[Unit]
Description=ProjectPlanner Live Inbox poll (oneshot)
After=network-online.target projectplanner.service projectplanner-gateway.service
Wants=network-online.target

[Service]
Slice=projectplanner-batch.slice
Type=oneshot
Nice=10
IOSchedulingClass=idle
MemoryMax=256M
User=ubuntu
WorkingDirectory=/opt/projectplanner
EnvironmentFile=/opt/projectplanner/.env
ExecStart=/opt/projectplanner/.venv/bin/python /opt/projectplanner/jobs.py poll_inbox
