# deploy/projectplanner-merge-coordinator.timer
# /etc/systemd/system/projectplanner-merge-coordinator.timer
# Runs the Switchboard merge-coordinator every 5 minutes (offset from claim-gate polling).
[Unit]
Description=ProjectPlanner Switchboard merge-coordinator (every 5 minutes)

[Timer]
OnBootSec=4min
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target
