# deploy/projectplanner-claim-gate.timer
# /etc/systemd/system/projectplanner-claim-gate.timer
# Poll open fleet PRs and refresh claim and merge-authorization statuses.
[Unit]
Description=ProjectPlanner Switchboard PR authorization gates (every 2 minutes)

[Timer]
OnBootSec=60s
OnUnitActiveSec=2min
Persistent=true

[Install]
WantedBy=timers.target
