# deploy/projectplanner-autodeploy.service
# /etc/systemd/system/projectplanner-autodeploy.service
# BUG-114: one-shot auto-deploy trigger. Fetches canonical master and, when the
# running checkout is behind, runs deploy/redeploy.sh (pull -> deps -> units ->
# restart -> Caddy fail-closed -> runtime proof + rollback). Refreshes the
# /health/version staleness signal on every tick. Triggered by the timer below.
#
# WHY THIS UNIT IS NOT SANDBOXED like the batch units (reconcile/narrate/etc.):
# its job is to mutate the running system — privileged git pull, systemctl
# restarts, and Caddy reload. It therefore runs as the ubuntu deploy login whose
# sudo policy already permits redeploy.sh. BUG-119 removes the reason for the
# temporary root workaround from BUG-114/PR #693: state and lock writes now land
# in narrow systemd-managed directories owned by ubuntu, never in the
# projectplanner-owned application data tree.
[Unit]
Description=ProjectPlanner Switchboard auto-deploy canonical master to prod
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
Nice=5
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/projectplanner
EnvironmentFile=/opt/projectplanner/.env
Environment=PLAN_ROOT=/opt/projectplanner
Environment=PYTHONDONTWRITEBYTECODE=1
# BUG-119: give the timer narrow, systemd-managed state and runtime directories.
# projectplanner.service reads the same 0644 signal file but cannot write the
# state directory. Neither path grants ubuntu access to the application DB tree.
# EnvironmentFile values override Environment values in systemd. Use a dedicated
# unit-owned variable so BUG-114's old PM_DEPLOY_STATE_FILE in .env cannot redirect
# this service back into the unwritable application data directory.
Environment=PM_DEPLOY_STATE_UNIT_FILE=/var/lib/projectplanner-autodeploy/deploy-state.json
Environment=AUTODEPLOY_LOCK_FILE=/run/projectplanner-autodeploy/auto-deploy.lock
StateDirectory=projectplanner-autodeploy
StateDirectoryMode=0755
RuntimeDirectory=projectplanner-autodeploy
RuntimeDirectoryMode=0700
# The script locks the whole cycle (atomic mkdir), so an overlapping tick or a
# tick landing on top of a long redeploy exits 0 without a second deploy.
ExecStart=/usr/bin/env bash /opt/projectplanner/deploy/auto_deploy.sh
# A failed redeploy has already rolled the edge back inside redeploy.sh; surface
# the non-zero result to the journal but never let systemd retry-loop a deploy.
TimeoutStartSec=1200
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
